# tasks/core/build_vm/overlay/etc/systemd/system/apt-cacher-ng.service.d/distro-seed.conf
[Unit]
Requires=ds-vm-runtime-setup.service
After=ds-vm-runtime-setup.service
